# night-suspend.service
[Unit]
Description=Suspend overnight

[Service]
Type=oneshot
ExecStart=sleepy-time.sh "12 hours"
User=root
Group=systemd-journal
