# suspend-over-weekend.timer
[Unit]
Description=Suspend over weekend

[Timer]
OnCalendar=Fri 20:30
Unit=weekend-suspend.service

[Install]
WantedBy=timers.target
